# examples/centos-bottle-mongo/supervisor.needrun.ini
[program:needrun]
user = needrun
command=/var/lib/needrun/env/bin/python /needrun/needrun/runserver.py
directory=/needrun/
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,PATH="/var/lib/needrun/env/bin/:/usr/bin:/bin"
stdout_logfile=/var/log/supervisor/needrun.log
redirect_stderr = true
autostart = true
